# build/service/coldfusion11/run
#!/bin/bash

set -e

exec /opt/coldfusion11/cfusion/bin/coldfusion start
