# build/service/coldfusion10/run
#!/bin/bash

set -e

exec /opt/coldfusion10/cfusion/bin/coldfusion start
